# 0 ressources/2023/pgl/crontab/root
0 4 * * * /opt/back/backup_incr_log.sh
30 4 * * * /opt/back/backup_incr_data.sh
